# recipes-core/udev/udev-rules-st/10-bcmdhd.rules
SUBSYSTEM=="module", ACTION=="add", RUN+="/sbin/insmod /lib/firmware/brcm/bcmdhd"
