# recipes-app/myd-ya157c/files/01-brcm.rules
SUBSYSTEM=="module", ACTION=="add", RUN+="/sbin/insmod /usr/share/100ask/bcmdhd.ko"